# debian/postrm
#! /bin/sh
# postrm script for file-rc
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postrm> `remove'
#        * <postrm> `purge'
#        * <old-postrm> `upgrade' <new-version>
#        * <new-postrm> `failed-upgrade' <old-version>
#        * <new-postrm> `abort-install'
#        * <new-postrm> `abort-install' <old-version>
#        * <new-postrm> `abort-upgrade' <old-version>
#        * <disappearer's-postrm> `disappear' <r>overwrit>r> <new-version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

# see #812337 + #285835, #805487
undivert() {
  dpkg-divert --package file-rc --remove --rename \
    --divert /usr/sbin/update-rc.d.init-system-helpers /usr/sbin/update-rc.d

  dpkg-divert --package file-rc --remove --rename \
    --divert /usr/sbin/invoke-rc.d.init-system-helpers /usr/sbin/invoke-rc.d

  dpkg-divert --package file-rc --remove --rename \
    --divert /usr/share/man/man8/update-rc.d.8.gz.init-system-helpers /usr/share/man/man8/update-rc.d.8.gz

  dpkg-divert --package file-rc --remove --rename \
    --divert /usr/share/man/man8/invoke-rc.d.8.gz.init-system-helpers /usr/share/man/man8/invoke-rc.d.8.gz
}

case "$1" in
    upgrade|failed-upgrade)
    ;;

    abort-upgrade)
      if dpkg --compare-versions "$2" lt 0.8.17 ; then
        undivert
      fi
    ;;

    abort-install|disappear)
        undivert
    ;;

    remove)
        undivert

	if [ ! -e /usr/sbin/update-rc.d ] ; then
		ln -sf /usr/share/sysvinit/update-rc.d /usr/sbin/update-rc.d
	fi
	if [ ! -e /usr/sbin/invoke-rc.d ] ; then
		ln -sf /usr/share/sysvinit/update-rc.d /usr/sbin/invoke-rc.d
	fi
    ;;

    purge)
	rm -f /var/backups/rc-links.tar.gz \
	    /etc/runlevel.conf /etc/runlevel.fallback
	echo Please purge /var/backups/rc.* manually.
    ;;

    *)
        echo "postrm called with unknown argument \`$1'" >&2
        exit 1
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
