# debian/preinst
#! /bin/sh
# preinst script for file-rc
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <new-preinst> `install'
#        * <new-preinst> `install' <old-version>
#        * <new-preinst> `upgrade' <old-version>
#        * <old-preinst> `abort-upgrade' <new-version>
#
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    install)
    ;;

    upgrade)
	#
	# Correct defects in runlevel.{conf|fallback} from file-rc (<< 0.4.1):
	#
	for i in runlevel.fallback runlevel.conf; do
	    sed -e 's/\([0-9][0-9]	-\)[	]\+\(\/etc\/init.d\)/\1	-		\2/g' \
		/etc/$i > /etc/$i.new
	    if ! cmp -s /etc/$i /etc/$i.new; then
		mv /etc/$i.new /etc/$i
		echo Corrected problem in /etc/$i.
	    else
		rm -f /etc/$i.new
	    fi
	done

	#
	# Undo diversions from file-rc (<<0.8.0):
	#
	if [ "$2" != "" ] && dpkg --compare-versions $2 lt 0.8.0
	then
	    if [ -f /usr/sbin/update-rc.d.links ]
	    then
		mv /usr/sbin/update-rc.d.links /usr/sbin/update-rc.d
	    fi
	    dpkg-divert --package file-rc --remove \
		--divert /usr/sbin/update-rc.d.links /usr/sbin/update-rc.d
	    if [ -f /usr/sbin/invoke-rc.d.links ]
	    then
		mv /usr/sbin/invoke-rc.d.links /usr/sbin/invoke-rc.d
	    fi
	    dpkg-divert --package file-rc --remove \
		--divert /usr/sbin/invoke-rc.d.links /usr/sbin/invoke-rc.d
	    if [ -f /etc/init.d/rc.links ]
	    then
		mv -f /etc/init.d/rc.links /etc/init.d/rc
	    fi
	    dpkg-divert --package file-rc --remove \
		--divert /etc/init.d/rc.links /etc/init.d/rc
	    if [ -f /etc/init.d/rcS.links ]
	    then
		mv -f /etc/init.d/rcS.links /etc/init.d/rcS
	    fi
	    dpkg-divert --package file-rc --remove \
		--divert /etc/init.d/rcS.links /etc/init.d/rcS
	fi
    ;;

    abort-upgrade)
    ;;

    *)
        echo "preinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
